# deployment/systemd/textreaderrpi.service
[Unit]
Description=TextReaderRpi - Display Text Driver Service
Documentation=file:///opt/textreaderrpi/docs/deployment/production-guide.md
After=network.target
Wants=network.target
[Service]
Type=simple
User=textreaderrpi
Group=textreaderrpi
# Runtime environment
WorkingDirectory=/opt/textreaderrpi
EnvironmentFile=-/opt/textreaderrpi/.env
ExecStart=/usr/bin/java -jar /opt/textreaderrpi/TextReaderRpi-all.jar
# Restart policy (D-33: on-failure with backoff)
Restart=on-failure
RestartSec=5s
RestartPreventExitStatus=0
# Health gate: wait for service to be ready before marking active
ExecStartPost=/bin/bash -c 'for i in $(seq 1 12); do sleep 5; curl -sf http://localhost:8080/health/ready && exit 0; done; exit 1'
# Resource limits
LimitNOFILE=65536
TimeoutStartSec=90s
TimeoutStopSec=30s
# Security hardening (T-03-09)
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/textreaderrpi/logs
# Standard I/O
StandardOutput=journal
StandardError=journal
SyslogIdentifier=textreaderrpi
[Install]
WantedBy=multi-user.target
